# salt-master.conf
[program:salt-master]
user=root
directory=/etc/salt
command=/usr/bin/salt-master -l debug
redirect_stderr=true

