# salt-api.conf
[program:salt-api]
user=root
directory=/etc/salt
command=/usr/bin/salt-api -l debug
redirect_stderr=true
autorestart=true