# salt-minion.conf
[program:salt-minion]
user=root
directory=/etc/salt
command=/usr/bin/salt-minion -l debug
redirect_stderr=true
